# home/dot_config/autostart/emacs-start.desktop
[Desktop Entry]
Exec=systemctl --user restart emacs
Icon=emacs
Name=emacs-start
Type=Application
