# home/dot_config/autostart/ssh-agent-start.desktop
[Desktop Entry]
TryExec=ssh-agent
Exec=systemctl --user restart ssh-agent
Name=Start ssh-agent
Type=Application
