# home/dot_config/profile.d/profile.sh
for name in gui linux_tty misc path xdg; do
  envsrc $name
done
