# home/dot_config/systemd/user/clash-gen.service
[Unit]
Description=Clash Config Generator

[Service]
Type=simple
ExecStart=/usr/bin/env ${HOME}/.local/bin/scripts/clash-gen
