# home/dot_config/systemd/user/clash-gen.service
[Unit]
Description=Clash Config Generator

[Service]
Type=simple
ExecStart=/bin/env clash-gen
