# home/dot_config/systemd/user/clash-localhost-gen.service
[Unit]
Description=Clash Config Generator for localhost

[Service]
Type=simple
ExecStart=/bin/env clash-gen -d ${HOME}/.config/clash-localhost
