# home/dot_config/systemd/user/envvar-setting.service
# -*- mode: conf -*-

[Unit]
Description=Environmental Variables Setting

[Service]
Type=oneshot
ExecStart=/usr/bin/true
RemainAfterExit=yes
