# home/dot_config/systemd/user/network-wait-online.service
# -*- mode: conf -*-

[Unit]
Description=network wait online

[Service]
Type=oneshot
ExecStart=nm-online -s -q
RemainAfterExit=yes
