# home/dot_config/systemd/user/pacman-mirrorlist-gen.service.tmpl
[Unit]
Description=Pacman Mirrorlist Generator

[Service]
Type=simple
ExecStart=env ${HOME}/.local/bin/scripts/pacman-mirrorlist-gen
