# home/dot_config/systemd/user/synloc-coor.service.tmpl
# -*- mode: conf -*-

[Unit]
Description=geographical coordinates update

[Service]
Type=simple
ExecStart=/usr/bin/env ${HOME}/.local/bin/scripts/synloc-coor
